# osl40562.service
[Unit]
Description=osl40562

[Service]
ExecStart=/usr/local/bin/osl40562

[Install]
WantedBy=default.target


